# owncloud/owncloud-supervisor.conf
[program:cron]
command = /usr/sbin/cron -f

[program:phpfpm]
command = /owncloud-run.sh
