# phabricator-aphlict/aphlict-supervisor.ini
[program:aphlict]
command = node --max-stack-size=256 /opt/phabricator/support/aphlict/server/aphlict_server.js --config=/opt/phabricator/conf/aphlict/aphlict.default.json
exitcodes = 0
startsecs = 0
user = phd
