# phabricator-phpfpm/phpfpm-supervisor.conf
[program:phpfpm]
command = /run.sh

[program:daemons]
command = /phabricator/phabricator/bin/phd start
exitcodes = 0
startsecs = 0
