# phabricator-ssh/sshd-supervisor.ini
[program:sshd]
command = /usr/sbin/sshd -D -f /etc/sshd/sshd_config

[program:syslogd]
command = rsyslogd -n
