# phabricator-ssh/sshd-supervisor.ini
[program:sshd]
command = /usr/sbin/sshd -D -E /var/log/sshd -f /etc/sshd/sshd_config

[program:keygen]
command = ssh-keygen -A
exitcodes = 0
startsecs = 0
autorestart = false
